# scripts/claude_polling_launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude.tradingagent.polling</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/santiagodiaz/Documents/Claude/Projects/Agent With API Key, For Trading Technical Analysis/scripts/run_claude_polling.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/santiagodiaz/Documents/Claude/Projects/Agent With API Key, For Trading Technical Analysis</string>
  <key>StandardOutPath</key>
  <string>/Users/santiagodiaz/Documents/Claude/Projects/Agent With API Key, For Trading Technical Analysis/launchd_claude_polling.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/santiagodiaz/Documents/Claude/Projects/Agent With API Key, For Trading Technical Analysis/launchd_claude_polling.log</string>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
