# scripts/discord_bot_launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude.tradingagent.discordbot</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/santiagodiaz/code/trading-agent/scripts/run_discord_bot.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/santiagodiaz/code/trading-agent</string>
  <key>StandardOutPath</key>
  <string>/Users/santiagodiaz/code/trading-agent/launchd_discord_bot.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/santiagodiaz/code/trading-agent/launchd_discord_bot.log</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
